# openWRT/openwrt.init.script
#!/bin/sh /etc/rc.common

START=30
STOP=80


start() {
        echo "Stopping previous oor process ..."
        killall /usr/sbin/oor &> /dev/null
        rm /var/run/oor.pid
        echo "Starting Open Overlay Router ..."
        /usr/sbin/oor -D
}

stop() {
        echo "Stopping Open Overlay Router ..."
        killall oor
}
